# workspace/etc/ld.so.preload
/home/junior/Desktop/Coding/agents_meta_hackathon/workspace/usr/local/lib/libchaos.so
