# workspace/etc/systemd/system/demo-api.service
[Unit]
Description=Demo API Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/app
EnvironmentFile=/etc/demo-api/demo-api.env
ExecStart=/usr/bin/python3 /opt/app/app.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
