# workspace/etc/cron.d/healthcheck
*/5 * * * * root /usr/local/bin/healthcheck.sh
